# visibli/host_conf/backup_user.crontab
@daily ssh -i /home/chris/Documents/keys/aws/server1.pem -q -C -p 22 %s rdiff-backup --server" ubuntu@XXXXX.amazonaws.com::/home/ubuntu/data/visibli_url_grab/ /home/chris/Documents/programming/visibli_grab/ec2_data/
