# visibli/host_conf/user.crontab
@daily sqlite3 /home/ubuntu/data/visibli_url_grab/visibli.db .dump > /home/ubuntu/data/visibli_url_grab/visibli.db.dump.sql
