# host_conf/user.crontab
@daily sqlite3 DATABASE_PATH .dump > DUMP_LOCATION
